# script.service
[Unit]

Description=ss

After=default.target

[Service]

ExecStart=/home/admin/ss.sh

[Install]

WantedBy=default.target
