# script.service
[Unit]

Description=ss

After=default.target

[Service]

ExecStart=ss.sh

[Install]

WantedBy=default.target
