# script.service
[Unit]

Description=ss

After=default.target

[Service]

ExecStart=/PATH/TO/ss.sh

[Install]

WantedBy=default.target
